# charan/.supervisor_conf.d/code-server.conf
[program:code-server]
command=bash -c 'code-server --bind-addr 0.0.0.0:8080 /app/Project'
autostart=true
autorestart=true
stderr_logfile=/var/log/code-server.err.log
priority=999
