# charan/.supervisor_conf.d/self-ping.conf
[program:Self-ping]
command=/usr/bin/python3 /app/self-ping.py
autostart=true
autorestart=true
stderr_logfile=/var/log/code-server.err.log
priority=6
